# .crontab
22 9 15 * * cd /home/ppham/mail; split-mail.pl | mail -s "Split-mail" ppham
