# crontab
PATH=/usr/local/sbin:/usr/local/bin:/usr/bin:/bin
MAILTO=ppham@local-box.org
* * * * * . /Users/ppham/venv/bs4/bin/activate; python3 /Users/ppham/src/blockchain/bitmain-scraper/scraper.py
